# systemd/pentarou.service
[Unit]
Description=Pentarou — Infrastructure Update Announcer
After=network.target

[Service]
Type=simple
User=pentarou
EnvironmentFile=/opt/pentarou/.env
WorkingDirectory=/opt/pentarou
ExecStart=/opt/pentarou/pentarou -config /opt/pentarou/config/config.yml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
