# systemd/pentarou-warn.service
[Unit]
Description=Pentarou pre-update warning
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/pentarou
ExecStart=/opt/pentarou/.venv/bin/python -m src warn
Environment=PENTAROU_MATRIX_TOKEN=
EnvironmentFile=-/opt/pentarou/.env
